# guardian/service/com.openclaw.guardian.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.guardian</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/oc-guardian</string>
        <string>--config</string>
        <string>/usr/local/etc/oc-guardian/guardian.toml</string>
        <string>start</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Only restart on abnormal exit (crash), not on clean exit (oc-guardian stop) -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Minimum 30s between restarts to prevent rapid restart loops -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>WorkingDirectory</key>
    <string>/usr/local/etc/oc-guardian</string>

    <key>StandardOutPath</key>
    <string>/usr/local/var/log/oc-guardian/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/oc-guardian/stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
